# systemd-services/kdeconnectd.service
[Service]
User=%i
ExecStart=/usr/lib/kdeconnectd
