# systemd-services/suspend.service
[Unit]
Description=Suspend

[Service]
Type=simple
ExecStart=/usr/bin/systemctl suspend
